# cookbooks/edb/templates/edb-as-9.6.service.erb
# It's not recommended to modify this file in-place, because it will be
# overwritten during package upgrades.  If you want to customize, the
# best way is to create a file "/etc/systemd/system/edb-as-9.6.service",
# containing
#       .include /lib/systemd/system/edb-as-9.6.service
#       ...make your changes here...
# For more info about custom unit files, see
# http://fedoraproject.org/wiki/Systemd#How_do_I_customize_a_unit_file.2F_add_a$

# Note: changing PGDATA will typically require adjusting SELinux
# configuration as well.

[Unit]
Description=EDB Postgres Advanced Server 9.6
After=syslog.target network.target

[Service]
Type=notify

User=enterprisedb
Group=enterprisedb

# Location of database directory
Environment=PGDATA=/opt/PostgresPlus/9.6AS/data
PIDFile=/var/lib/edb/as9.6/data/postmaster.pid

# Where to send early-startup messages from the server (before the logging
# options of postgresql.conf take effect)
# This is normally controlled by the global default set by systemd
# StandardOutput=syslog

# Disable OOM kill on the postmaster
OOMScoreAdjust=-1000

ExecStartPre=/usr/edb/as9.6/bin/edb-as-96-check-db-dir ${PGDATA}
ExecStart=/usr/edb/as9.6/bin/edb-postmaster -D ${PGDATA}
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT

# Give a reasonable amount of time for the server to start up/shut down
TimeoutSec=300

[Install]
WantedBy=multi-user.target

